# roles/supervisor-mailhog/templates/mailhog.conf.j2
[program:mailhog]
command=/usr/local/bin/mailhog
process_name=%(program_name)s
startsecs=0
user=root
stdout_logfile=/var/log/mailhog/stdout.log
stderr_logfile=/var/log/mailhog/stderr.log